# conf/roomsrv.conf
[program:roomsvr]
priority=20
directory=/root/apprtc
command=/root/google_appengine/dev_appserver.py --host=0.0.0.0 /root/apprtc/out/app_engine
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log














